# units/vyatta-static-arp.service.d/danos.conf
[Service]
# Only one ExecStart option is allowed, so reset it first
ExecStart=
ExecStart=/opt/vyatta/sbin/arp.py /etc/vyatta-routing/frr.json
