# units/frr.service.d/danos.conf
[Unit]
PartOf=vyatta-routing-frr.target
Before=vyatta-routing-frr-early.target
Wants=vyatta-frr-base-config.service
