# units/vyatta-routing-frr-early.target
[Unit]
Description=Vyatta Routing Service (early)
Before=vyatta-routing-frr.target
RefuseManualStart=true
RefuseManualStop=true
